# projects/recaria/recaria.service
# Systemd service file for recaria
# /etc/systemd/system/recaria.service

[Unit]
Description=Recaria Django Application
After=network.target

[Service]
Type=notify
User=ubuntu
Group=ubuntu
WorkingDirectory=/home/ubuntu/recaria/recaria
Environment=PATH=/usr/local/bin:/usr/bin:/bin
Environment=PYTHONPATH=/home/ubuntu/recaria
Environment=DJANGO_SETTINGS_MODULE=recaria_backend.settings
ExecStart=/usr/local/bin/gunicorn --config /home/ubuntu/recaria/gunicorn.conf.py recaria_backend.wsgi:application
ExecReload=/bin/kill -s HUP $MAINPID
KillMode=mixed
TimeoutStopSec=5
PrivateTmp=true
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target

